# deploy/enjazy.service
[Unit]
Description=Enjazy Portfolio App
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/enjazy
Environment=PORT=8000
ExecStart=/opt/enjazy/.venv/bin/python /opt/enjazy/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
